# deploy/workbench.service
[Unit]
Description=Workbench Operations System (Hunchentoot/Common Lisp)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/workbench
ExecStart=/usr/bin/sbcl --load /root/quicklisp/setup.lisp --load /opt/workbench/start.lisp --eval "(loop (sleep 3600))"
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
